# examples/hydra_node_beta.service
[Unit]
Description=NDN Hydra Service (Node Beta)

[Service]
Type=simple
ExecStart=/home/hydra_op/ndn-hydra/venv/bin/ndn-hydra-repo -rp hydra-test@bruins -n beta

[Install]
WantedBy=multi-user.target
~
